# content/service/vaultwarden/run
#!/bin/sh

cd /

exec 2>&1
exec /start.sh